# tarsnap/files/tarsnap.timer
[Unit]
Description=run tarsnap daily and on startup

[Timer]
OnActiveSec=30min
OnCalendar=daily

[Install]
WantedBy=timers.target
